# unicorn.conf
[program:unicorn]
directory=/thanks_ants
command=bundle exec unicorn -E production -c config/unicorn.rb
stdout_events_enabled=true
stderr_events_enabled=true
